# api.service
[Unit]
Description=API server

[Service]
Alias=api.service
ExecStart=/tmp/api

